# containers/www/supervisord.conf
[supervisord]
nodaemon=true

[program:tornado]
directory=/www
command=/usr/bin/python3 tornado-ws.py
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
